# templates/default/kafka-upstart.conf.erb
# Start the Kafka Broker daemon
#
author          "Mike Heffner <mike@librato.com>"
description     "Runs the Kafka Broker daemon"

start on runlevel [2345]
stop on runlevel [!2345]

# XXX: force the right Java
env JAVA_HOME=<%= @java_home %>

respawn

exec <%= @kafka_dir %>/bin/kafka-server-start.sh <%= @kafka_dir %>/config/server.properties >> /var/log/kafka/broker.log 2>&1
